# cron.d/fix-stale-handle
PATH="/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"
*/5 * * * * root /usr/local/bin/fix-stale-handle cron
